# mkosi.extra-iscsi/usr/lib/systemd/system/iscsi-shutdown.service.d/disable-in-initrd.conf
# Make the unit do nothing in the initrd.
# We can't mask the unit because iscsi.services has Requires=iscsi-shutdown.service,
# so instead we let it survive switch-root.
[Unit]
IgnoreOnIsolate=true

# This is a bug in the upstream unit file. This unit has DefaultDependencies=no,
# and Conflicts=shutdown.target, but no ordering, which only guarantees that the
# stop job will be queued.
Before=shutdown.target
